# docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[program:redis]
command=redis-server --maxmemory 256mb --maxmemory-policy allkeys-lru --save "" --daemonize no
priority=10
autorestart=true
startretries=5
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:api]
command=uvicorn api.main:app --host 127.0.0.1 --port 8000 --log-level info
directory=/app
priority=20
autorestart=true
startretries=3
startsecs=3
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
environment=REDIS_HOST="127.0.0.1",REDIS_PORT="6379",PYTHONPATH="/app",PYTHONUNBUFFERED="1"

[program:worker]
command=python -m worker.runner
directory=/app
priority=30
autorestart=true
startretries=3
startsecs=5
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
environment=REDIS_HOST="127.0.0.1",REDIS_PORT="6379",PYTHONPATH="/app",PYTHONUNBUFFERED="1"

[program:nginx]
command=nginx -g "daemon off;"
priority=40
autorestart=true
startretries=3
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

